# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:node]
command=/bin/bash -c "cd /language-server && npm run start"
